# deployment/daily-checkin.timer
[Unit]
Description=Run Daily Check-in Automation at 22:22

[Timer]
OnCalendar=*-*-* 22:22:00
Persistent=true
Unit=daily-checkin.service

[Install]
WantedBy=timers.target
